# root/etc/s6-overlay/s6-rc.d/rclone/finish
#!/command/execlineb -S0
# shellcheck shell=bash

/run/s6/basedir/bin/halt
